# build/shareclipext.conf
[program:shareclipext]
user = nobody
command = /usr/bin/shareclipext server 3000
umask = 000

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh '^/usr/bin/shareclipext'
umask = 000